# supervisord-play.conf
[program:play]
command=/start-play.sh
numprocs=1
autostart=true
autorestart=true
